# examples/yubikey@.service
[Service]
Type=forking
ExecStart=/usr/local/bin/yubilock systemd-event /home/you/yubilock.yaml %I
